# app/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:tabulae-celery]
command=/apps/tabulae1/bin/celery --app=tabulae.celery:app worker --loglevel=INFO --concurrency=1
directory=/apps/tabulae-v2/app ;
user=ubuntu
environment=TABULAE_ENVIRONMENT='prod',
numprocs=1
stdout_logfile=/apps/log/celery-worker.log
stderr_logfile=/apps/log/celery-worker.log
autostart=true
autorestart=true
startsecs=10

[program:flower]
user=ubuntu
command=/apps/tabulae1/bin/celery flower --app=tabulae.celery:app worker --address=0.0.0.0 --port=5555 --basic_auth=newsai:cUj2KAH6NKNoGjFDp ;
directory=/apps/tabulae-v2/app ;
environment=TABULAE_ENVIRONMENT='prod',
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/apps/log/out_celeryd_flower.log
sterr_logfile=/apps/log/error_celeryd_flower.log

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998


[inet_http_server]
port=127.0.0.1:9001   ;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001 ;
